# com.sourdough.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sourdough.tunnel</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/cloudflared</string>
        <string>tunnel</string>
        <string>--url</string>
        <string>http://localhost:8080</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/moltbot/Projects/sourdough-monitor/data/tunnel_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/moltbot/Projects/sourdough-monitor/data/tunnel_stderr.log</string>

    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
